# board/alpinelinux/overlay/etc/init.d/calamares
#!/sbin/openrc-run
#
#  SPDX-License-Identifier: GPL-2.0+
#
#  This file is part of linux-distros-br2-external.
#
#  linux-distros-br2-external is free software; you can redistribute it and/or
#  modify it under the terms of the GNU General Public License as published by
#  the Free Software Foundation; either version 2.0 of the License, or (at your
#  option) any later version.

depend() {
	use dbus
	after dbus
}

start() {
	ebegin "Calamares Installer"
	if calamares -platform eglfs; then
		rm -f /etc/runlevels/*/$RC_SVCNAME
	fi
	eend $?
}
